# valkyrie_desktop/etc/skel/.config/autostart/Dropbox.desktop
[Desktop Entry]
Exec=~/.dropbox-dist/dropboxd &
Name=Dropbox
Type=Application
Version=1.0
X-LXQt-Need-Tray=true
