# valkyrie_desktop/etc/skel/.config/autostart/dunst.desktop
[Desktop Entry]
Exec=dunst
Name=dunst
Type=Application
Version=1.0
