# valkyrie_desktop/etc/skel/.config/autostart/network-manager.desktop
[Desktop Entry]
Exec=nm-applet
Name=network-manager
Type=Application
Version=1.0
X-LXQt-Need-Tray=true
