# valkyrie_desktop/etc/skel/.config/autostart/openbox_autostart.desktop
[Desktop Entry]
Exec=~/.config/openbox/autostart
Name=openbox_autostart
Type=Application
Version=1.0
