# valkyrie_desktop/etc/skel/.config/autostart/radeon-profile.desktop
[Desktop Entry]
Exec=/usr/bin/radeon-profile
Name=radeon-profile
Type=Application
Version=1.0
X-LXQt-Need-Tray=true
